# rootfs/etc/s6-overlay/s6-rc.d/msmtpd/run
#!/command/with-contenv sh
# Wait for msmtp.conf to be generated by signpost-web
echo "msmtpd: waiting for /data/signpost/msmtp.conf..."
while [ ! -f /data/signpost/msmtp.conf ]; do
  sleep 2
done
echo "msmtpd: starting on 127.0.0.1:2500"
exec msmtpd --interface=127.0.0.1 --port=2500 --log=/tmp/msmtp.log --command="msmtp -C /data/signpost/msmtp.conf --read-envelope-from --"
