# rootfs/etc/s6-overlay/s6-rc.d/maddy/run
#!/command/with-contenv sh
sleep 3
exec /bin/maddy -config /data/signpost/maddy.conf run
