# rootfs/etc/s6-overlay/s6-rc.d/signpost-web/run
#!/command/with-contenv sh
exec /app/signpost
